# common/roles/dm-ai.logstash/templates/supervisor.k8s-log.ini
[program:k8s-logstash]
environment = JAVA_HOME="/usr/lib/jvm/java-1.8.0-openjdk-1.8.0.222.b10-0.el7_6.x86_64"
command= {{ logstash.install_path }}/{{ logstash.tar_dir_name }}/bin/logstash -f {{ logstash.install_path }}/{{ logstash.tar_dir_name }}/config/logstash.k8s-log.conf --path.data {{ logstash.install_path }}/k8s/data
process_name=%(program_name)s-%(process_num)d
numprocs=1
autostart=true
autorestart=true
user=root
stdout_logfile=/data/logs/supervisor/logstash_k8s_stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/data/logs/supervisor/logstash_k8s_stderr.log
stderr_logfile_maxbytes=10MB
