# common/roles/prometheus/node_exporter_systemctl/templates/node_exporter.service.j2
[Service]
User=root
Group=root
ExecStart=/usr/bin/dmai_node_exporter --web.listen-address=:{{D_port}}

[Install]
WantedBy=multi-user.target

[Unit]
Description=node_exporter
After=network.target
