# common/roles/prometheus/node_exporter/templates/supervisor.node-exporter.ini
[program:{{ project_name }}]
command=/usr/bin/{{ bin_name }} --web.listen-address=":{{ D_port }}"
process_name=%(program_name)s-%(process_num)d
numprocs=1
autostart=true
autorestart=true
user=root
stdout_logfile=/data/logs/supervisor/{{ project_name }}_stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/data/logs/supervisor/{{ project_name }}_stderr.log
stderr_logfile_maxbytes=10MB
